# ex04/33-ngbanza-usb.rules
SUBSYSTEM=="input", SUBSYSTEMS=="usb", ACTION=="add", RUN+="/sbin/insmod /root/home/Little-penguin/ex04/main.ko"
